# deploy/weather-trip-scout.service
[Unit]
Description=Weather Trip Scout morning report

[Service]
Type=oneshot
WorkingDirectory=/opt/weather-trip-scout
ExecStart=/opt/weather-trip-scout/.venv/bin/python -m app.main run
EnvironmentFile=/opt/weather-trip-scout/.env
