# deploy/weather-trip-scout.timer
[Unit]
Description=Run weather-trip-scout every morning

[Timer]
Unit=weather-trip-scout.service
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
